# modules/aws_asg/unregister-netdevice.crontab
SHELL = /bin/bash

* * * * * root /var/tmp/travis-run.d/handle-unregister-netdevice &>/dev/null
