# modules/aws_asg_fixed/unregister-netdevice.crontab
SHELL = /bin/bash

* * * * * root /var/tmp/travis-run.d/check-unregister-netdevice &>/dev/null
