# assets/travis-worker/clean-up-containers.crontab
SHELL = /bin/bash

* * * * * root /var/tmp/travis-run.d/clean-up-containers &>> /var/log/syslog
